# config/drc_tc_gui.conf.in
# drc_tc - MIT DARPA Robotics Challenge Traffic Controlled Tunnel GUI
#
# GUI for controlling drc_tc simulated restricted network tunnel

description	"MIT DARPA Robotics Challenge Traffic Controlled Tunnel GUI"

env GOBY_LIAISON_PLUGINS=@CMAKE_INSTALL_PREFIX@/lib/libliaison_plugins_drc_tc.so

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

script
   /usr/bin/goby_liaison /etc/drc_tc/drc_goby_liaison.pb.cfg
end script
