# config/drc_tc_link2.conf.in
# drc_tc - MIT DARPA Robotics Challenge Traffic Controlled Tunnel
#
# Simulates restricted access tunnel

description	"MIT DARPA Robotics Challenge Traffic Controlled Tunnel"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

pre-start script
    set -e
    echo "Pre-start"
    test -x /usr/bin/socat || { stop; exit 0; }
    test -x @CMAKE_INSTALL_PREFIX@/bin/drc_blackout || { stop; exit 0; }
    @CMAKE_INSTALL_PREFIX@/bin/drc_tc_parse /etc/drc_tc/drc_tc.pb.cfg > /tmp/drc_tc_env.sh
    chmod u+x /tmp/drc_tc_env.sh

   @CMAKE_INSTALL_PREFIX@/bin/drc_tc_parse /etc/drc_tc/drc_tc.pb.cfg > /tmp/drc_tc_env.sh
   chmod u+x /tmp/drc_tc_env.sh
   cat >> /tmp/drc_tc_env.sh << EOF
   
    /usr/bin/socat TUN:\$DRC_TC_LINK2_TUNNEL_ADDRESS/24,up,tun-name=tun\$DRC_TC_LINK2_TUNNEL_NUM,iff-no-pi exec:"@CMAKE_INSTALL_PREFIX@/bin/drc_blackout \$DRC_TC_LINK2_REMOTE_ADDRESS \$DRC_TC_LINK2_REMOTE_PORT \$DRC_TC_LINK2_BLACKOUT_FILE",nofork

EOF
end script

exec /tmp/drc_tc_env.sh

post-start script
   set -e
    echo "Post-start"
end script

post-stop script
end script
